# clipthread-tkinter-ui.service
[Unit]
Description=ClipThread Tkinter UI Service
After=network.target

[Service]
Type=simple
User=%i
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
Environment=DATABASE_URL=/home/%i/.clipthread/db.sqlite
ExecStart=/home/%i/.local/bin/clipthread-ui
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target